# sshserver-sshd.service
[Unit]
Description=SSH TUI Portfolio Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tuiserver
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 